# docker/runit/scheduler/run
#!/bin/sh

exec setuser www-data /etc/service/scheduler/scheduler-rake.sh
