# docker/runit/couchdb/run
#!/bin/sh

mkdir -p /var/run/couchdb /data/couchdb
chown -R couchdb:couchdb /var/run/couchdb /etc/couchdb /data/couchdb
exec setuser couchdb /usr/bin/couchdb
