# docker/runit/couchdb/run
#!/bin/sh

export HOME=/var/lib/couchdb
mkdir -p /var/run/couchdb
chown -R couchdb:couchdb /var/run/couchdb
chown -R couchdb:couchdb /etc/couchdb/local.d/*.ini
exec chpst -u couchdb /usr/bin/couchdb
