# docker/runit/worker/run
#!/bin/bash

sv -v -w 15 check /etc/service/couchdb up
sv -v -w 15 check /etc/service/solr up
echo "Starting Delayed Job Worker in $RAILS_ENV mode..."
cd /rapidftr
bundle exec rake jobs:work
