# lib/tasks/ubuntu/firstboot/rapidftr.conf
# Upstart script for starting RapidFTR Solr and Scheduler

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [016]

post-start script
  su -l ubuntu -c 'cd /srv/rapid_ftr_production/current && RAILS_ENV=production bundle exec rake sunspot:solr:restart'
  su -l ubuntu -c 'cd /srv/rapid_ftr_production/current && RAILS_ENV=production bundle exec rake scheduler:restart'
end script

pre-stop script
  su -l ubuntu -c 'cd /srv/rapid_ftr_production/current && RAILS_ENV=production bundle exec rake sunspot:solr:stop'
  su -l ubuntu -c 'cd /srv/rapid_ftr_production/current && RAILS_ENV=production bundle exec rake scheduler:stop'
end script
