# lib/tasks/ubuntu/rapidftr.conf
# Upstart script for starting RapidFTR Solr and Scheduler

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [016]

post-start script
  lxc-start -n rapidftr_lxc -d
  lxc-wait  -n rapidftr_lxc -s RUNNING

  RAPIDFTR_LXC_HOST="10.0.3.3"
  (
    iptables -t nat -A PREROUTING -p tcp --dport 80    -j DNAT --to-destination $RAPIDFTR_LXC_HOST:80
    iptables -t nat -A PREROUTING -p tcp --dport 443   -j DNAT --to-destination $RAPIDFTR_LXC_HOST:443
    iptables -t nat -A PREROUTING -p tcp --dport 5984  -j DNAT --to-destination $RAPIDFTR_LXC_HOST:5984
    iptables -t nat -A PREROUTING -p tcp --dport 6984  -j DNAT --to-destination $RAPIDFTR_LXC_HOST:6984
    iptables -A FORWARD -p tcp -o lxcbr0 -j ACCEPT
  ) || echo "Unable to configure port forwarding rules"

  echo "RapidFTR started"
end script

post-stop script
  RAPIDFTR_LXC_HOST="10.0.3.3"
  (
    iptables -t nat -D PREROUTING -p tcp --dport 80    -j DNAT --to-destination $RAPIDFTR_LXC_HOST:80
    iptables -t nat -D PREROUTING -p tcp --dport 443   -j DNAT --to-destination $RAPIDFTR_LXC_HOST:443
    iptables -t nat -D PREROUTING -p tcp --dport 5984  -j DNAT --to-destination $RAPIDFTR_LXC_HOST:5984
    iptables -t nat -D PREROUTING -p tcp --dport 6984  -j DNAT --to-destination $RAPIDFTR_LXC_HOST:6984
    iptables -D FORWARD -p tcp -o lxcbr0 -j ACCEPT
  ) || echo "Unable to remove port forwarding rules"

  lxc-stop -n rapidftr_lxc
  lxc-wait -n rapidftr_lxc -s STOPPED
  echo "RapidFTR stopped"
end script
